# wildlife-loop.service
[Unit]
Description=Wildlife Autonomous Tuning Loop — deterministic pre-gate + Claude judgment session
After=network-online.target
Wants=network-online.target

[Service]
# Type=oneshot: each tick is a single invocation that exits when done.
# Paired with wildlife-loop.timer (narrowed to rough night hours as defence-in-depth).
#
# TICK FLOW (three stages):
#
# 1. GATE (loop.nightgate) — ~0.5s, zero tokens:
#    Checks two conditions deterministically (no LLM):
#      a. Is it currently night? (SunChecker.is_daytime() == False — astral, Bonn)
#      b. Has this loop-day already completed?
#         (state.json["last_tick_completed_day"] == loop_day() — stable across midnight)
#    On SKIP: sends a heartbeat Telegram ping once per loop-day (deadman guard), then
#    exits 1. The `|| exit 0` turns this into a clean unit exit (no failed-unit noise).
#    On PROCEED: exits 0, allowing the shell && chain to continue to stage 2.
#
# 2. JUDGMENT (claude -p) — Opus session, consumes the loop.md prompt:
#    Runs the full nightly analysis: ingest → adjudicate → metrics → decide → deploy
#    → report. Runs only when nightgate exits 0.
#
# 3. COMPLETION STAMP (loop.endtick) — stamped by the LLM as its FINAL step
#    (see experiments/loop.md), NOT by systemd. Rationale: tying "tonight is done"
#    to claude's exit code is unsafe — claude -p can exit 0 on a clean-but-incomplete
#    run (early stop, or a usage-limit truncation), which would falsely mark the night
#    done and skip the remaining work. Instead the LLM runs `python -m loop.endtick`
#    only after it actually finishes the workflow, stamping
#    state.json["last_tick_completed_day"] = loop_day(). Any interruption before that
#    point (limit, crash, hang) leaves it unstamped, so the next tick RESUMES from
#    committed state. Failure mode is benign (re-run), never "skip needed work".
#
# PYTHONPATH/uv note: uv run resolves the venv correctly; PYTHONPATH=src is set
# via Environment= so `-m loop.*` finds the src/loop package without a pip install.
#
# Validated headless invocation (no TTY required):
#   systemd-run --user --pipe --wait \
#     bash -c 'claude -p --dangerously-skip-permissions --model opus < experiments/loop.md'
# confirmed exit 0 with a trivial probe prompt, 2026-06-08, claude 2.1.168, Pi 5 aarch64.
# Authentication: OAuth via ~/.claude/.credentials.json (no API key needed).
# Note: systemd's StandardInput=file: does NOT work with claude (it times out waiting
# on stdin); the prompt must be fed via a shell stdin redirect as below.
Type=oneshot
User=daniel
WorkingDirectory=/home/daniel/animal_tracker
Environment="PATH=/home/daniel/.local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/daniel/animal_tracker/src"
# nightgate: skip → `|| exit 0` (unit succeeds; heartbeat already handled inside nightgate)
#            proceed → claude runs the full tick and stamps its own completion via
#            `loop.endtick` as its final step (see stage 3 above — deliberately NOT a
#            systemd `&& endtick`, so a clean-but-incomplete claude exit can't false-mark
#            the night done).
# --dangerously-skip-permissions: required for unattended operation.
# --model opus: the judgment layer requires the stronger model (as specified in loop.md).
ExecStart=/bin/bash -c 'cd /home/daniel/animal_tracker && PYTHONPATH=src /home/daniel/.local/bin/uv run python -m loop.nightgate || exit 0; /home/daniel/.local/bin/claude -p --dangerously-skip-permissions --model opus < /home/daniel/animal_tracker/experiments/loop.md'
StandardOutput=journal
StandardError=journal

# Resource limits
Nice=10

[Install]
WantedBy=multi-user.target
