# wildlife-loop.service
[Unit]
Description=Wildlife Autonomous Tuning Loop — per-tick Claude judgment session
After=network-online.target
Wants=network-online.target

[Service]
# Type=oneshot: each tick is a single claude -p invocation that exits when done.
# Paired with wildlife-loop.timer (OnCalendar=every 2h, Persistent=true).
# The loop.md prompt is self-gating: it checks state.json and stops if it is not
# night or tonight's run is already complete, so a fresh session per tick is both
# correct and safe.
#
# Validated headless invocation (no TTY required):
#   systemd-run --user --pipe --wait \
#     bash -c 'claude -p --dangerously-skip-permissions --model opus < experiments/loop.md'
# confirmed exit 0 with a trivial probe prompt, 2026-06-08, claude 2.1.168, Pi 5 aarch64.
# Authentication: OAuth via ~/.claude/.credentials.json (no API key needed).
# Note: systemd's StandardInput=file: does NOT work with claude (it times out waiting
# on stdin); the prompt must be fed via a shell stdin redirect as below.
Type=oneshot
User=daniel
WorkingDirectory=/home/daniel/animal_tracker
Environment="PATH=/home/daniel/.local/bin:/usr/bin:/bin"
# --dangerously-skip-permissions: required for unattended operation (no human to approve tool calls).
# --model opus: the judgment layer requires the stronger model (as specified in loop.md).
# The prompt is fed over stdin via the `< experiments/loop.md` redirect (not as a shell
# argument), so edits take effect on the next tick AND the prompt content can never break
# the unit's shell quoting (e.g. an apostrophe in loop.md). Only the fixed command and
# file path live in the bash -c string.
ExecStart=/bin/bash -c '/home/daniel/.local/bin/claude -p --dangerously-skip-permissions --model opus < /home/daniel/animal_tracker/experiments/loop.md'
StandardOutput=journal
StandardError=journal

# Resource limits
Nice=10

[Install]
WantedBy=multi-user.target
