# wildlife-loop.service
[Unit]
Description=Wildlife Autonomous Tuning Loop — deterministic pre-gate + Claude judgment session
After=network-online.target
Wants=network-online.target

[Service]
# Type=oneshot: each tick is a single invocation that exits when done.
# Paired with wildlife-loop.timer (narrowed to rough night hours as defence-in-depth).
#
# DETERMINISTIC PRE-GATE (loop.nightgate):
#   Before launching the expensive Opus session, a cheap Python module checks two
#   conditions without any LLM call (~0.5s, zero tokens):
#     1. Is it currently night? (SunChecker.is_daytime() == False, same logic as
#        the camera main loop — both use astral with Bonn coordinates from Config)
#     2. Has tonight's run already completed? (state.json last_metrics.date != today)
#   exit 0 = loop should run; exit 1 = skip this tick.
#   The `|| exit 0` ensures a skip (exit 1 from the gate) makes the unit exit 0,
#   so skipped ticks never show as failed units in `systemctl status`.
#
# PYTHONPATH/uv note: uv run resolves the venv correctly; PYTHONPATH=src is set
# via Environment= (same pattern as wildlife-deploy.service) so `-m loop.nightgate`
# finds the src/loop package without a pip install step.
#
# Validated headless invocation (no TTY required):
#   systemd-run --user --pipe --wait \
#     bash -c 'claude -p --dangerously-skip-permissions --model opus < experiments/loop.md'
# confirmed exit 0 with a trivial probe prompt, 2026-06-08, claude 2.1.168, Pi 5 aarch64.
# Authentication: OAuth via ~/.claude/.credentials.json (no API key needed).
# Note: systemd's StandardInput=file: does NOT work with claude (it times out waiting
# on stdin); the prompt must be fed via a shell stdin redirect as below.
Type=oneshot
User=daniel
WorkingDirectory=/home/daniel/animal_tracker
Environment="PATH=/home/daniel/.local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/daniel/animal_tracker/src"
# Gate first: if nightgate exits 1 (skip), the `|| exit 0` short-circuits and the
# unit exits 0 cleanly. Only if the gate exits 0 (proceed) does claude run.
# --dangerously-skip-permissions: required for unattended operation.
# --model opus: the judgment layer requires the stronger model (as specified in loop.md).
ExecStart=/bin/bash -c 'cd /home/daniel/animal_tracker && PYTHONPATH=src /home/daniel/.local/bin/uv run python -m loop.nightgate || exit 0; /home/daniel/.local/bin/claude -p --dangerously-skip-permissions --model opus < /home/daniel/animal_tracker/experiments/loop.md'
StandardOutput=journal
StandardError=journal

# Resource limits
Nice=10

[Install]
WantedBy=multi-user.target
